# templates/default/jenkins.service.erb
[Unit]
Description=Jenkins
After=network.target
Requires=network.target

[Service]
Type=simple
EnvironmentFile=<%= node['airgapped_jenkins']['master']['environment_file_path'] %>/jenkins
ExecStart=/usr/bin/java <%= node['airgapped_jenkins']['master']['jvm_options'] %> -jar <%= node['airgapped_jenkins']['master']['home'] %>/jenkins.war --httpPort=<%= node['airgapped_jenkins']['master']['port'] %> --httpListenAddress=<%= node['airgapped_jenkins']['master']['listen_address'] %> <%= node['airgapped_jenkins']['master']['jenkins_args'] %>
Restart=always
User=<%= node['airgapped_jenkins']['master']['user'] %>
RestartSec=20

[Install]
WantedBy=multi-user.target