# templates/default/jenkins-slave-init.erb
#!/bin/sh
# chkconfig: 345 99 01
# description: Manages a Jenkins JNLP slave service

#
# This file was generated by Chef for <%= node['fqdn'] %>.
# Do NOT modify this file by hand.
#

JENKINS_USER=<%= @new_resource.user %>
export HOME=<%= @new_resource.remote_fs %>
export JENKINS_HOME=$HOME
logfile=$JENKINS_HOME/$(basename $0).log

cmd="<%= @java_bin %> -jar <%= @slave_jar %> -jnlpUrl <%= @jnlp_url %>"
<%- if @jvm_options -%>
cmd="$cmd <%= @jvm_options %>"
<%- end -%>
<%- if @jnlp_secret -%>
cmd="$cmd -secret <%= @jnlp_secret %>"
<%- end -%>

do_start()
{
  cd $JENKINS_HOME
  su $JENKINS_USER -c "$cmd >>$logfile 2>&1" &
}

do_stop()
{
  for pid in $(ps -eo pid,cmd | awk '/-jnlpUrl <%= Regexp.escape(@jnlp_url).gsub('/', '\/') %>/ && !/awk/ {print $1}'); do
    kill $pid
  done
}

case "$1" in
  start)
    do_start
    ;;
  stop)
    do_stop
    ;;
  restart)
    do_stop
    sleep 3
    do_start
    ;;
esac
